# tools/package/postinst
#!/bin/bash

echo "Post installation script"

#Update sysctl configuration
sysctl --system

#Enable Service
systemctl enable poco-boot.service
